# .agents/configs/aidevops-profile-readme-update.service
# systemd user service for profile README auto-update (Linux equivalent of launchd plist)
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp aidevops-profile-readme-update.service ~/.config/systemd/user/
#   cp aidevops-profile-readme-update.timer ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now aidevops-profile-readme-update.timer
#
# Check status:
#   systemctl --user status aidevops-profile-readme-update.timer
#   journalctl --user -u aidevops-profile-readme-update.service

[Unit]
Description=aidevops profile README stats update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.aidevops/agents/scripts/profile-readme-helper.sh update
StandardOutput=append:%h/.aidevops/.agent-workspace/logs/profile-readme-update.log
StandardError=append:%h/.aidevops/.agent-workspace/logs/profile-readme-update.log

[Install]
WantedBy=default.target
