# .agents/configs/sh.aidevops.profile-readme-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>sh.aidevops.profile-readme-update</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>log="$HOME/.aidevops/.agent-workspace/logs/profile-readme-update.log"; mkdir -p "$(dirname "$log")"; printf 'Profile README launchd wrapper started at %s\n' "$(date -u +%Y-%m-%dT%H:%M:%SZ)" >> "$log"; exec "$HOME/.aidevops/agents/scripts/profile-readme-helper.sh" update >> "$log" 2&gt;&amp;1</string>
	</array>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardOutPath</key>
	<string>/dev/null</string>
	<key>StandardErrorPath</key>
	<string>/dev/null</string>
</dict>
</plist>
