# .agents/templates/launchd/com.aidevops.aidevops-supervisor-merge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- SPDX-License-Identifier: MIT -->
<!-- SPDX-FileCopyrightText: 2025-2026 Marcus Quinn -->
<!--
  Template: com.aidevops.aidevops-supervisor-merge.plist (t21247, GH#21247)

  Decoupled merge-pass launchd agent. Runs pulse-wrapper.sh with the
  merge-only flag every 60 seconds so green PRs merge within ~60-90 s of
  CI completion, independent of the main dispatch cycle (~23 min average).

  The actual installed plist at ~/Library/LaunchAgents/ is generated by
  setup_pulse_merge_routine() in setup-modules/schedulers-platform.sh with
  real absolute paths substituted. This template is the canonical shape.

  Key differences from com.aidevops.aidevops-supervisor-pulse.plist:
    StartInterval = 60s (not the configurable pulse dispatch interval).
    ProgramArguments adds the merge-only flag.
    StandardOutPath/StandardErrorPath route to pulse-merge.log.
    Label: com.aidevops.aidevops-supervisor-merge
    No KeepAlive (merge pass exits cleanly).
    RunAtLoad false (first run at the 60s tick).
    ProcessType Background + LowPriorityBackgroundIO + Nice 10.

  Lock: ~/.aidevops/locks/pulse-merge-instance.lock (separate lock).
  Log:  ~/.aidevops/logs/pulse-merge.log

  See also: .agents/scripts/pulse-wrapper.sh
            .agents/scripts/pulse-wrapper-bootstrap.sh
            setup-modules/schedulers-platform.sh
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.aidevops.aidevops-supervisor-merge</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/bash</string>
		<string>/Users/YOU/.aidevops/agents/scripts/pulse-wrapper.sh</string>
		<string>--merge-only</string>
	</array>
	<key>StartInterval</key>
	<integer>60</integer>
	<key>StandardOutPath</key>
	<string>/Users/YOU/.aidevops/logs/pulse-merge.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOU/.aidevops/logs/pulse-merge.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>/Users/YOU</string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>LowPriorityBackgroundIO</key>
	<true/>
	<key>Nice</key>
	<integer>10</integer>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>4096</integer>
	</dict>
</dict>
</plist>
