# .agents/templates/launchd/com.aidevops.aidevops-supervisor-merge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- SPDX-License-Identifier: MIT -->
<!-- SPDX-FileCopyrightText: 2025-2026 Marcus Quinn -->
<!--
  Template: com.aidevops.aidevops-supervisor-merge.plist (t21247, GH#21247)

  Decoupled merge-pass launchd agent — runs pulse-wrapper.sh --merge-only every
  60 seconds so green PRs merge within ~60-90 s of CI completion, independent of
  the main dispatch cycle (~23 min average).

  The actual installed plist at ~/Library/LaunchAgents/ is generated by
  setup_pulse_merge_routine() in setup-modules/schedulers.sh with real absolute
  paths substituted. This template file is the canonical reference shape.

  Key differences from com.aidevops.aidevops-supervisor-pulse.plist:
    - StartInterval = 60 (not the configurable pulse dispatch interval ~120-180 s)
    - ProgramArguments adds --merge-only flag
    - StandardOutPath/StandardErrorPath → pulse-merge.log (not pulse-wrapper.log)
    - Label: com.aidevops.aidevops-supervisor-merge
    - No OPENCODE_BIN, PULSE_DIR, PULSE_STALE_THRESHOLD (not used by merge-only path)
    - No KeepAlive (merge pass exits cleanly; no restart-on-crash needed)
    - RunAtLoad = false (first run at the 60s tick, not immediately at login)
    - ProcessType = Background + LowPriorityBackgroundIO + Nice = 10 (low-impact)

  Lock: ~/.aidevops/locks/pulse-merge-instance.lock (separate from main pulse lock)
  Log:  ~/.aidevops/logs/pulse-merge.log

  See also: .agents/scripts/pulse-wrapper.sh (--merge-only implementation)
            setup-modules/schedulers.sh (setup_pulse_merge_routine)
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.aidevops.aidevops-supervisor-merge</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/bash</string>
		<string>/Users/YOU/.aidevops/agents/scripts/pulse-wrapper.sh</string>
		<string>--merge-only</string>
	</array>
	<key>StartInterval</key>
	<integer>60</integer>
	<key>StandardOutPath</key>
	<string>/Users/YOU/.aidevops/logs/pulse-merge.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOU/.aidevops/logs/pulse-merge.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>/Users/YOU</string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>LowPriorityBackgroundIO</key>
	<true/>
	<key>Nice</key>
	<integer>10</integer>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>4096</integer>
	</dict>
</dict>
</plist>
